# system-config/etc/systemd/system/vtrgb-nord.service
[Unit]
Description=Apply Nord VGA palette to virtual console
After=systemd-vconsole-setup.service
Before=greetd.service

[Service]
Type=oneshot
ExecStart=/usr/bin/setvtrgb /etc/vtrgb
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
